# system/upstart/express-app.service
[Unit]
After=network.target vagrant.mount

[Service]
ExecStart=/vagrant/start-server.sh

[Install]
WantedBy=default.target
